# turtle-ha/rootfs/etc/s6-overlay/s6-rc.d/turtle/run
#!/command/with-contenv bashio
exec /usr/bin/turtle-run.sh
